# mkosi.extra/etc/systemd/system/console-getty.service.d/10-no-vhangup.conf
[Service]
TTYVHangup=no
CapabilityBoundingSet=~CAP_SYS_TTY_CONFIG
