# omit-fp/mkosi.extra/etc/systemd/system/bench.service.d/execstart.conf
[Service]
ExecStart=/root/src/bench.py --omit-fp pyperformance
